# etc/services.d/logs/run
#!/bin/bash

exec outputlogs -o `cat stdout` -e `cat stderr`
